# deploy/systemd/marketplace-scrape-loop.service
[Unit]
Description=Marketplace Ouedkniss scrape-and-seed iteration
Documentation=file:/opt/marketplace/scripts/run-loop.sh
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
# --quiet suppresses the OK summary on stdout; everything goes to the
# per-run log + structured metrics (see run-loop.sh --help). systemd
# captures any unsuppressed output to journald automatically.
# --max-products caps the seller catalog by deleting oldest products
# after each seed (cascades media+variants).
# --categories rotates one Ouedkniss top-level category per run via
# the state file. Each category keeps its own next_start_page entry,
# so progress is preserved across rotations.
ExecStart=/opt/marketplace/scripts/run-loop.sh \
  --seller-id 019e08a4-97cd-7d98-afd7-670878dc51c2 \
  --categories telephones,informatique,electronique_electromenager,vetements_mode,sante_beaute,immobilier,automobiles_vehicules \
  --max-products 280000 \
  --quiet
Restart=no
TimeoutStartSec=240
User=root
Group=root
ProtectSystem=strict
ReadWritePaths=/opt/marketplace/data
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
